# files/infisical-cookbooks/package/templates/default/infisical-runsvdir.service.erb
[Unit]
Description=Infisical Runit supervision process
After=<%= @systemd_after %>

[Service]
ExecStart=/opt/infisical/embedded/bin/runsvdir-start
Restart=always
<% if @tasks_max %>
TasksMax=<%= @tasks_max %>
<% end %>

[Install]
WantedBy=<%= @systemd_wanted_by %>
